# scheduler.cron
54 18 * * * root . /etc/environment && cd /app && /usr/local/bin/python manage.py run_daily_pipeline >> /var/log/cron.log 2>&1
